# deploy/com.clawd.cc-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clawd.cc-bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/austingriffith/clawd/clawd-harness/projects/clawd-video-chat/cc-bridge.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/austingriffith/clawd/clawd-harness/projects/clawd-video-chat</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>CC_BRIDGE_MODEL</key>
        <string>opus</string>
        <key>CLAUDE_P_AGENT_HOME</key>
        <string>/Users/austingriffith/clawd/clawd-harness/projects/claude-p-agent</string>
        <key>CC_BRIDGE_CWD</key>
        <string>/Users/austingriffith/clawd/clawd-harness/projects/claude-p-agent</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/austingriffith/.cache/clawd/cc-bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/austingriffith/.cache/clawd/cc-bridge.log</string>
</dict>
</plist>
